# deploy/glimmer.service
[Unit]
Description=Glimmer
Documentation=https://github.com/
After=network.target
Wants=network.target

[Service]
Type=simple
User=glimmer
Group=glimmer
WorkingDirectory=/opt/glimmer
ExecStart=/opt/glimmer/glimmer -config /etc/glimmer/config.toml
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glimmer

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/glimmer

[Install]
WantedBy=multi-user.target
